# profile.d/system-wide.sh
#
# System-wide bash settings
#

alias ll='ls -FGlAh --color=auto --group-directories-first --indicator-style=file-type'
alias ls='ls --color=auto'
